# Resources/redisdata.conf
[program:redisdata]
command=redis-server /var/www/go/src/thyself/data
directory=/var/www/go/src/thyself/data
environment=PATH="/var/www/go/src/thyself/data"
user=ubuntu
