# Resources/thygo.conf
[program:thygo]
command=/var/www/go/src/thyself/thyself
directory=/var/www/go/src/thyself
environment=PATH="/var/www/go/src/thyself"
user=ubuntu
